# config/systemd/whisper-gui.service
[Unit]
Description=Whisper Voice Recording GUI - faster-whisper speech-to-text with system tray
After=graphical-session-started.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=WHISPER_PROJECT_DIR
ExecStart=WHISPER_PROJECT_DIR/venv/bin/python3 -m whisper_app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="PATH=WHISPER_PROJECT_DIR/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=WHISPER_PROJECT_DIR/src"
Environment="PYTHONUNBUFFERED=1"
Environment="QT_QPA_PLATFORM=xcb"
Environment="QT_QPA_SYSTEMTRAYICON=xcb"

[Install]
WantedBy=graphical-session.target
